# crontabs/root
* */12 * * * /usr/share/nginx/html/scripts/index.sh
